# deb-package/lib/systemd/system/config-pilot.service
[Unit]
Description=Config Pilot - GitOps Configuration Monitor
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/config-pilot
Environment=CONFIG_PATH=/etc/config-pilot/config.yaml
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
